# usr/local/lib/systemd/system/systemd-homed-zfs-udev-mount@.service
[Unit]
Description=systemd-homed ZFS mounting helper script

[Service]
ExecStart=/usr/local/bin/my-services/zfs-homed-udev-mount.sh %i
